# init.d/rfsd.kamikaze
#!/bin/sh /etc/rc.common

START=50

start()
{
	rfsd
}

stop()
{
	killall rfsd
}
