# init.d/rfsd.kamikaze
#!/bin/sh /etc/rc.common

START=50

LISTEN="0.0.0.0"
PORT="5001"

start()
{
	rfsd -a ${LISTEN} -p ${PORT}
}

stop()
{
	killall rfsd
}
